# packages/graphics/hyperion/system.d/hyperion.service
[Unit]
Description=Hyperion Daemon
After=graphical.target

[Service]
ExecStartPre=/usr/bin/gpio2spi
ExecStart=/bin/sh /usr/local/hyperion-wrapper
